# faq-mcp.service
[Unit]
Description=FAQ MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Kshitij/Bootcamp-mcp
Environment="PATH=/home/ubuntu/Kshitij/Bootcamp-mcp/venv/bin"
EnvironmentFile=/home/ubuntu/Kshitij/Bootcamp-mcp/.env
ExecStart=/home/ubuntu/Kshitij/Bootcamp-mcp/venv/bin/python mcp-server/faq_server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=faq-mcp

[Install]
WantedBy=multi-user.target
